# ops/quadlet/adhocint-web.container
[Unit]
Description=adhocint web
After=adhocint-db.service adhocint-minio.service
Requires=adhocint-db.service adhocint-minio.service

[Container]
Image=localhost/adhocint-web:latest
ContainerName=adhocint-web
EnvironmentFile=/opt/adhocint/shared/.env.production
Environment=HOSTNAME=0.0.0.0
PublishPort=127.0.0.1:3001:3001
Network=adhocint

[Install]
WantedBy=default.target
